# debian/postinst
#!/bin/bash
#

systemctl daemon-reload
systemctl unmask galliumos-kabylake
systemctl enable galliumos-kabylake
systemctl start galliumos-kabylake

exit 0
